# etc/voucher.timer
[Unit]
Description=Run Voucher Service Daily at Midnight

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target